# media_server/tasks/templates/media-sync.service.j2
[Unit]
Description=Sync media from central NAS
Wants=media-sync.timer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sync_from_nas.sh

