# media_server/tasks/templates/media-sync.timer.j2
[Unit]
Description=Nightly media sync

[Timer]
OnCalendar=*-*-* {{ sync_time }}
Persistent=true

[Install]
WantedBy=timers.target

